# scripts/com.vetka.uds-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vetka.uds-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/danilagulin/Documents/VETKA_Project/vetka_live_03/.venv/bin/python3</string>
        <string>/Users/danilagulin/Documents/VETKA_Project/vetka_live_03/scripts/uds_daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/vetka-uds-daemon.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/vetka-uds-daemon.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/danilagulin/Documents/VETKA_Project/vetka_live_03</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/danilagulin/Documents/VETKA_Project/vetka_live_03</string>
    </dict>
</dict>
</plist>
